# recipes-connectivity/hsic-hubs/files/hsic-hubs.service
[Unit]
Description=HSICHubs
After=network.target

[Service]
ExecStart=/bin/sh /usr/lib/hsic-hubs/service.sh

[Install]
WantedBy=multi-user.target
